# system-scripts/etc/systemd/system/mapbuilder.service
[Unit]
Description=Map Builder

[Service]
ExecStart=/opt/mapbuilder
Restart=always

[Install]
WantedBy=multi-user.target